# roles/common/files/autostart/guake.desktop
[Desktop Entry]
Type=Application
Name=Guake
Exec=/usr/bin/guake
Terminal=false
OnlyShowIn=GNOME;
