# roles/common/files/autostart/keepass.desktop
[Desktop Entry]
Type=Application
Name=Keepass
Exec=/usr/bin/keepass
Terminal=false
NoDisplay=true
